# deploy/pitch-deck-analyzer.service
[Unit]
Description=Pitch Deck Analyzer API
After=network.target

[Service]
Type=simple
User=ec2-user
WorkingDirectory=/home/ec2-user/pitch-deck-analyzer
ExecStart=/home/ec2-user/pitch-deck-analyzer/venv/bin/python /home/ec2-user/pitch-deck-analyzer/api_server.py
Restart=on-failure
RestartSec=5
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
